# systemd/seavpn-monitor.service
[Unit]
Description=SeaVPN Service Monitor
After=network.target
Wants=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/root
ExecStart=/usr/bin/python3 /root/service_monitor.py
Restart=always
RestartSec=10
StandardOutput=syslog
StandardError=syslog
SyslogIdentifier=seavpn-monitor

[Install]
WantedBy=multi-user.target
