# systemd/seavpn-admin.service
[Unit]
Description=SeaVPN Admin Web Panel
After=network.target
Wants=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/root
ExecStart=/usr/bin/python3 /root/admin_web.py
Restart=always
RestartSec=10
StandardOutput=syslog
StandardError=syslog
SyslogIdentifier=seavpn-admin

[Install]
WantedBy=multi-user.target


