# systemd/seavpn-webhook.service
[Unit]
Description=SeaVPN Webhook Handler
After=network.target
Wants=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/root
Environment="YOOKASSA_SHOPID=1141477"
Environment="YOOKASSA_SECRET_KEY=live_QMWqzDoOaDHGhRW5vlfU9nNM8rRdjStxe1pUayD_30w"
Environment="BOT_TOKEN=8394863099:AAGQ3w7Y7WnYr9kfji_UcFc3jV-I--TIIvQ"
ExecStart=/usr/bin/python3 /root/webhook_handler.py
Restart=always
RestartSec=10
StandardOutput=syslog
StandardError=syslog
SyslogIdentifier=seavpn-webhook

[Install]
WantedBy=multi-user.target
